# configs/ubuntu/onboot.service
[Unit]
Description=Runs on boot
After=multi-user.target

[Service]
Type=idle
ExecStart=sh /etc/onboot

[Install]
WantedBy=multi-user.target