# configs/ubuntu/xmodmapinit.conf
description "xmodmapinit"
start on desktop-start
stop on desktop-end

script 
    if [ -s ~/.xmodmap ]; then
        xmodmap ~/.xmodmap
    fi
end script